# package-scripts/contegix-puppet/postinst
#!/bin/bash
#
# Perform necessary puppet setup steps
# after package is installed.
#


prefix="contegix-"
puppet_services="puppet-agent " 
install_base="/opt/contegix/puppet"
if [ -f /etc/redhat-release ]; then
  	os_family="RedHat"
	os_major_rev="$(perl -pe 's/.*release (\d+)\.\d+.*/$1/' /etc/redhat-release)"
else
  	os_family="Unknown"
	os_major_rev="Unknown"
fi

PROGNAME=$(basename $0)

function error_exit
{
  echo "${PROGNAME}: ${1:-"Unknown Error"}" 1>&2
  exit 1
}

for service in $puppet_services; do
	if [ -e /etc/init.d/${prefix}$service ] ; then 
		chkconfig ${prefix}$service off
                chkconfig --del ${prefix}${service}
                rm -f /etc/init.d/${prefix}${service}
	fi
done

# link logrotate configuration
[ -h /etc/logrotate.d/${prefix}puppet ] || ln -s /opt/contegix/puppet/etc/logrotate.d/puppet /etc/logrotate.d/${prefix}puppet
# link hourly puppet agent cron
[ -h /etc/cron.hourly/${prefix}puppet-agent ] || ln -s ${install_base}/etc/cron.hourly/puppet-agent /etc/cron.hourly/${prefix}puppet-agent

exit 0
