# package-scripts/contegix-puppet/postinst
#!/bin/bash
#
# Perform necessary puppet setup steps
# after package is installed.
#


prefix="contegix-"
puppet_services="puppet-agent " 
install_base="/opt/contegix/puppet"
PROGNAME=$(basename $0)
OS_MAJOR_REV="$( ${install_base}/bin/facter operatingsystemmajrelease)"
OS_FAMILY="$( ${install_base}/bin/facter osfamily )"

function error_exit
{
  echo "${PROGNAME}: ${1:-"Unknown Error"}" 1>&2
  exit 1
}

for service in $puppet_services; do
	if [ ! -e /etc/init.d/${prefix}$service ] ; then 
		ln -s /opt/contegix/puppet/etc/init.d/$service /etc/init.d/${prefix}${service}
		chkconfig --add ${prefix}$service
		chkconfig ${prefix}$service off
	fi
done

# link logrotate configuration
[ -h /etc/logrotate.d/${prefix}puppet ] || ln -s /opt/contegix/puppet/etc/logrotate.d/puppet /etc/logrotate.d/${prefix}puppet

if [ "${OS_MAJOR_REV}" == "7" -a "${OS_FAMILY}" == "RedHat" ] ; then	
	[ -h /etc/systemd/system/${prefix}puppet-agent.service ] || ln -s ${install_base}/etc/systemd/system/contegix-puppet-agent.service /etc/systemd/system/${prefix}puppet-agent.service	
	systemctl daemon-reload
	systemctl disable contegix-puppet-agent.service
else
	# link hourly puppet agent cron
	[ -h /etc/cron.hourly/${prefix}puppet-agent ] || ln -s ${install_base}/etc/cron.hourly/puppet-agent /etc/cron.hourly/${prefix}puppet-agent
fi

exit 0
